# hermes-sidecar/launchd/io.solvys.fintheon-hermes.plist
<?xml version="1.0" encoding="UTF-8"?>
<!-- [claude-code 2026-04-19] Local launchd unit for the Hermes sidecar. Load with:
     launchctl load -w ~/Library/LaunchAgents/io.solvys.fintheon-hermes.plist
     Or symlink from the repo: ln -s $(pwd)/hermes-sidecar/launchd/io.solvys.fintheon-hermes.plist ~/Library/LaunchAgents/
-->
<!DOCTYPE plist PUBLIC "-//Apple//DTD PLIST 1.0//EN" "http://www.apple.com/DTDs/PropertyList-1.0.dtd">
<plist version="1.0">
<dict>
    <key>Label</key>
    <string>io.solvys.fintheon-hermes</string>
    <key>ProgramArguments</key>
    <array>
        <string>/bin/bash</string>
        <string>-lc</string>
        <string>cd "$HERMES_SIDECAR_DIR" &amp;&amp; exec "$UV_BIN" run python entrypoint.py</string>
    </array>
    <key>EnvironmentVariables</key>
    <dict>
        <key>PATH</key>
        <string>/opt/homebrew/bin:/usr/local/bin:/usr/bin:/bin</string>
        <key>HERMES_SIDECAR_DIR</key>
        <string>/Users/tifos/Documents/Codebases/fintheon/hermes-sidecar</string>
        <key>UV_BIN</key>
        <string>/opt/homebrew/bin/uv</string>
        <key>PORT</key>
        <string>8318</string>
        <key>HERMES_LOG_LEVEL</key>
        <string>info</string>
    </dict>
    <key>RunAtLoad</key>
    <true/>
    <key>KeepAlive</key>
    <dict>
        <key>SuccessfulExit</key>
        <false/>
        <key>Crashed</key>
        <true/>
    </dict>
    <key>ThrottleInterval</key>
    <integer>10</integer>
    <key>StandardOutPath</key>
    <string>/tmp/fintheon-hermes.log</string>
    <key>StandardErrorPath</key>
    <string>/tmp/fintheon-hermes.err.log</string>
</dict>
</plist>
