# backend-hono/scripts/com.fintheon.dispatch-mdb.plist
<?xml version="1.0" encoding="UTF-8"?>
<!DOCTYPE plist PUBLIC "-//Apple//DTD PLIST 1.0//EN" "http://www.apple.com/DTDs/PropertyList-1.0.dtd">
<plist version="1.0">
<dict>
    <key>Label</key>
    <string>com.fintheon.dispatch-mdb</string>

    <key>ProgramArguments</key>
    <array>
        <string>/bin/bash</string>
        <string>/Users/tifos/Desktop/Codebases/fintheon/backend-hono/scripts/dispatch-mdb.sh</string>
    </array>

    <!-- 6:30 AM ET, Mon-Fri -->
    <key>StartCalendarInterval</key>
    <array>
        <dict><key>Weekday</key><integer>1</integer><key>Hour</key><integer>6</integer><key>Minute</key><integer>30</integer></dict>
        <dict><key>Weekday</key><integer>2</integer><key>Hour</key><integer>6</integer><key>Minute</key><integer>30</integer></dict>
        <dict><key>Weekday</key><integer>3</integer><key>Hour</key><integer>6</integer><key>Minute</key><integer>30</integer></dict>
        <dict><key>Weekday</key><integer>4</integer><key>Hour</key><integer>6</integer><key>Minute</key><integer>30</integer></dict>
        <dict><key>Weekday</key><integer>5</integer><key>Hour</key><integer>6</integer><key>Minute</key><integer>30</integer></dict>
    </array>

    <key>StandardOutPath</key>
    <string>/Users/tifos/Desktop/Codebases/fintheon/backend-hono/logs/dispatch-mdb-launchd.log</string>
    <key>StandardErrorPath</key>
    <string>/Users/tifos/Desktop/Codebases/fintheon/backend-hono/logs/dispatch-mdb-launchd.log</string>

    <key>EnvironmentVariables</key>
    <dict>
        <key>PATH</key>
        <string>/usr/local/bin:/opt/homebrew/bin:/usr/bin:/bin:/Users/tifos/.bun/bin</string>
    </dict>

    <!-- Run missed job on wake/login if Mac was asleep at 6:30 -->
    <key>StartInterval</key>
    <integer>0</integer>
</dict>
</plist>
