# backend-hono/scripts/com.fintheon.dispatch-tott.plist
<?xml version="1.0" encoding="UTF-8"?>
<!DOCTYPE plist PUBLIC "-//Apple//DTD PLIST 1.0//EN" "http://www.apple.com/DTDs/PropertyList-1.0.dtd">
<plist version="1.0">
<dict>
    <key>Label</key>
    <string>com.fintheon.dispatch-tott</string>

    <key>ProgramArguments</key>
    <array>
        <string>/bin/bash</string>
        <string>/Users/tifos/Documents/Codebases/fintheon/backend-hono/scripts/dispatch-tott.sh</string>
    </array>

    <!-- 4:30 PM ET, Sunday only -->
    <key>StartCalendarInterval</key>
    <dict>
        <key>Weekday</key><integer>0</integer>
        <key>Hour</key><integer>16</integer>
        <key>Minute</key><integer>30</integer>
    </dict>

    <key>StandardOutPath</key>
    <string>/Users/tifos/Documents/Codebases/fintheon/backend-hono/logs/dispatch-tott-launchd.log</string>
    <key>StandardErrorPath</key>
    <string>/Users/tifos/Documents/Codebases/fintheon/backend-hono/logs/dispatch-tott-launchd.log</string>

    <key>EnvironmentVariables</key>
    <dict>
        <key>PATH</key>
        <string>/usr/local/bin:/opt/homebrew/bin:/usr/bin:/bin:/Users/tifos/.bun/bin</string>
    </dict>
</dict>
</plist>
