# launchd/io.solvys.fintheon-news-worker.plist
<?xml version="1.0" encoding="UTF-8"?>
<!DOCTYPE plist PUBLIC "-//Apple//DTD PLIST 1.0//EN" "http://www.apple.com/DTDs/PropertyList-1.0.dtd">
<plist version="1.0">
<dict>
  <key>Label</key>
  <string>io.solvys.fintheon-news-worker</string>
  <key>ProgramArguments</key>
  <array>
    <string>/Users/tifos/.bun/bin/bun</string>
    <string>run</string>
    <string>src/workers/riskflow-worker/index.ts</string>
  </array>
  <key>WorkingDirectory</key>
  <string>/Users/tifos/Documents/Codebases/fintheon/backend-hono</string>
  <key>KeepAlive</key>
  <dict>
    <key>SuccessfulExit</key>
    <false/>
    <key>Crashed</key>
    <true/>
  </dict>
  <key>RunAtLoad</key>
  <true/>
  <key>StandardOutPath</key>
  <string>/tmp/fintheon-news-worker.out.log</string>
  <key>StandardErrorPath</key>
  <string>/tmp/fintheon-news-worker.err.log</string>
  <key>EnvironmentVariables</key>
  <dict>
    <key>NEWS_WORKER_PORT</key>
    <string>8082</string>
    <key>FLAG_NEWS_WORKER_WRITES_RISKFLOW</key>
    <string>true</string>
    <key>PATH</key>
    <string>/Users/tifos/.bun/bin:/opt/homebrew/bin:/usr/local/bin:/usr/bin:/bin</string>
  </dict>
  <key>ThrottleInterval</key>
  <integer>10</integer>
</dict>
</plist>
